# lehr_03/configs/yannic.service
[Unit]
Description=Start Grafana and Prometheus from docker

[Service]
ExecStart=/etc/yannic/startup.sh

 

[Install]
WantedBy=multi-user.target